# raspi_user/pi-power-button-master/script.systemd/listen-for-shutdown.service
[Unit]
Description=pi-power-button

[Service]
User=root
WorkingDirectory=/root
ExecStart=/usr/local/bin/listen-for-shutdown.py
Restart=always

[Install]
WantedBy=multi-user.target

